# systemd/secrets-inject.service
[Unit]
Description=Inject pre-generated secrets from secrets ISO
# Must run before services that would generate keys
Before=age-generate-identity.service ssh-generate-identity.service
# Ensure disk subsystem is ready
After=local-fs.target
# Only run if secrets ISO is attached
ConditionPathExists=/dev/disk/by-label/SECRETS

[Service]
Type=oneshot
RemainAfterExit=yes

# Mount secrets ISO
ExecStartPre=/usr/bin/mkdir -p /var/mnt/secrets
ExecStartPre=/usr/bin/mount -o ro -t iso9660 /dev/disk/by-label/SECRETS /var/mnt/secrets

# Create target directories with secure permissions
ExecStartPre=/usr/bin/install -d -m 0700 /etc/age
ExecStartPre=/usr/bin/install -d -m 0700 /etc/git-ssh

# Copy age keys
ExecStartPre=/usr/bin/install -m 0400 /var/mnt/secrets/age.key /etc/age/identity.txt
ExecStartPre=/usr/bin/install -m 0444 /var/mnt/secrets/age.key.pub /etc/age/identity.pub

# Copy SSH keys
ExecStartPre=/usr/bin/install -m 0400 /var/mnt/secrets/ssh.key /etc/git-ssh/id_ed25519
ExecStartPre=/usr/bin/install -m 0444 /var/mnt/secrets/ssh.key.pub /etc/git-ssh/id_ed25519.pub

# Unmount secrets ISO
ExecStart=/usr/bin/umount /var/mnt/secrets
ExecStartPost=/usr/bin/rmdir /var/mnt/secrets

[Install]
WantedBy=multi-user.target
